# files/etc/uci-defaults/94-usb-packages-install
#!/bin/bash
# USB 软件包自动安装脚本
# 保存到 /etc/usb-packages-install.sh
# 在 /etc/rc.local 中添加：/etc/usb-packages-install.sh &

echo "=== USB 软件包自动安装脚本 ==="

# 检查 USB 存储是否挂载
USB_MOUNT="/mnt/usb"
PKG_DIR="$USB_MOUNT/packages"

# 尝试挂载 USB（如果未挂载）
if ! mountpoint -q "$USB_MOUNT"; then
    echo "尝试挂载 USB 存储..."
    mkdir -p "$USB_MOUNT"
    
    # 查找 USB 设备
    USB_DEVICE=$(blkid | grep -E 'TYPE="ext4"|TYPE="vfat"' | head -1 | cut -d: -f1)
    
    if [ -n "$USB_DEVICE" ]; then
        mount "$USB_DEVICE" "$USB_MOUNT"
        if [ $? -eq 0 ]; then
            echo "USB 挂载成功: $USB_DEVICE"
        else
            echo "USB 挂载失败"
            exit 1
        fi
    else
        echo "未找到 USB 存储设备"
        exit 1
    fi
fi

# 创建软件包目录
mkdir -p "$PKG_DIR"

# 修改 opkg 配置，使用 USB 作为存储
cat > /etc/opkg.conf << EOF
dest root /
dest ram /tmp
lists_dir ext /var/opkg-lists
option overlay_root /overlay
option check_signature
EOF

# 设置 opkg 缓存到 USB
sed -i "s|lists_dir ext /var/opkg-lists|lists_dir ext $PKG_DIR/opkg-lists|" /etc/opkg.conf

echo "=== 开始安装软件包到 USB ==="

# 安装 Docker（如果需要）
if [ ! -f /usr/bin/docker ]; then
    echo "安装 Docker..."
    opkg install --dest usb docker dockerd docker-compose luci-app-dockerman
fi

# 安装 Samba（如果需要）
if [ ! -f /usr/sbin/smbd ]; then
    echo "安装 Samba..."
    opkg install --dest usb samba4-server luci-app-samba4
fi

# 安装 AdGuardHome（如果需要）
if [ ! -f /usr/bin/AdGuardHome ]; then
    echo "安装 AdGuardHome..."
    opkg install --dest usb adguardhome luci-app-adguardhome
fi

# 创建符号链接，让系统能找到 USB 中的可执行文件
if [ -d "$PKG_DIR/usr/bin" ]; then
    for bin in $(ls "$PKG_DIR/usr/bin/"); do
        if [ ! -f "/usr/bin/$bin" ]; then
            ln -sf "$PKG_DIR/usr/bin/$bin" "/usr/bin/$bin"
        fi
    done
fi

if [ -d "$PKG_DIR/usr/sbin" ]; then
    for sbin in $(ls "$PKG_DIR/usr/sbin/"); do
        if [ ! -f "/usr/sbin/$sbin" ]; then
            ln -sf "$PKG_DIR/usr/sbin/$sbin" "/usr/sbin/$sbin"
        fi
    done
fi

echo "=== USB 软件包安装完成 ==="
echo "软件包安装在: $PKG_DIR"
echo "可用空间: $(df -h $USB_MOUNT | tail -1 | awk '{print $4}')"

# 重启相关服务
/etc/init.d/dockerd enable 2>/dev/null && /etc/init.d/dockerd start 2>/dev/null
/etc/init.d/samba4 enable 2>/dev/null && /etc/init.d/samba4 start 2>/dev/null
/etc/init.d/AdGuardHome enable 2>/dev/null && /etc/init.d/AdGuardHome start 2>/dev/null

echo "=== 安装完成，请重启路由器 ==="