# files/etc/uci-defaults/96-install-extras
#!/bin/sh
# 96-install-extras (原版 + V4 优化移植)

DONE_FLAG="/etc/install-extras-done"
[ -f "$DONE_FLAG" ] && exit 0

mkdir -p /usr/lib/install-extras

cat > /usr/lib/install-extras/run.sh << 'RUNEOF'
#!/bin/sh
LOG_TAG="install-extras"
DONE_FLAG="/etc/install-extras-done"
LOCK="/tmp/install-extras.lock"

log() { echo "[$(date '+%H:%M:%S')] $1" | logger -t "$LOG_TAG"; }

[ -f "$DONE_FLAG" ] && exit 0
[ -f "$LOCK" ] && exit 0
touch "$LOCK"

log "▶ install-extras 开始（优化版）"

# ── 清理 video feed ──────────────────────────────────────────────────
for f in /etc/apk/repositories.d/*.list; do
    [ -f "$f" ] && sed -i '/\/video\//d' "$f"
done

# ── 等待网络 ─────────────────────────────────────────────────────────
sleep 10
NETWORK_OK=0
for i in $(seq 1 60); do
    ping -c 1 -W 2 223.5.5.5 >/dev/null 2>&1 && { NETWORK_OK=1; break; }
    ping -c 1 -W 2 8.8.8.8  >/dev/null 2>&1 && { NETWORK_OK=1; break; }
    sleep 3
done
if [ "$NETWORK_OK" -eq 0 ]; then
    log "❌ 网络不可达，等下次 WAN UP 重试"
    rm -f "$LOCK"; exit 1
fi
log "✅ 网络就绪"

# ── 强制写入公共 DNS ─────────────────────────────────────────────────
printf 'nameserver 223.5.5.5\nnameserver 8.8.8.8\n' > /etc/resolv.conf
log "✅ DNS 注入：223.5.5.5 / 8.8.8.8"

# ── 检测 U 盘 ────────────────────────────────────────────────────────
USB_MOUNT=""
log "⏳ 检测 U 盘挂载..."
for i in $(seq 1 15); do
    for mp in /mnt/sda1 /mnt/sda /mnt/usb /mnt/mmcblk0p1; do
        mountpoint -q "$mp" 2>/dev/null && { USB_MOUNT="$mp"; break 2; }
    done
    sleep 2
done

# ═════════════════════════════════════════════════════════════════════
# 🔧 V4 优化点1：I/O 调度 & 内存优化
# ═════════════════════════════════════════════════════════════════════
if [ -n "$USB_MOUNT" ]; then
    log "✅ U 盘已挂载: $USB_MOUNT"
    
    # I/O 调度器优化
    USB_DEV=$(mount | grep "$USB_MOUNT" | awk '{print $1}' | sed 's/[0-9]*//g')
    if [ -n "$USB_DEV" ] && [ -f "/sys/block/$USB_DEV/queue/scheduler" ]; then
        if grep -q "mq-deadline" "/sys/block/$USB_DEV/queue/scheduler" 2>/dev/null; then
            echo "mq-deadline" > "/sys/block/$USB_DEV/queue/scheduler" 2>/dev/null
            log "⚡ I/O 调度器: mq-deadline"
        elif grep -q "deadline" "/sys/block/$USB_DEV/queue/scheduler" 2>/dev/null; then
            echo "deadline" > "/sys/block/$USB_DEV/queue/scheduler" 2>/dev/null
            log "⚡ I/O 调度器: deadline"
        fi
    fi
    
    # 降低 swappiness，减少 U 盘写入
    echo "10" > /proc/sys/vm/swappiness 2>/dev/null
    log "⚡ 内存优化: swappiness=10"
else
    log "⚠️ 未检测到 U 盘，docker/AdGuardHome 相关将跳过"
fi

# ═════════════════════════════════════════════════════════════════════
# 🔧 V4 优化点2：智能 Swap（自动创建）
# ═════════════════════════════════════════════════════════════════════
if [ -n "$USB_MOUNT" ] && [ ! -f "$USB_MOUNT/swapfile" ]; then
    MEM_TOTAL=$(free -k | awk '/^Mem:/ {print $2}')
    # 512MB 内存，创建 256MB swap 足够
    if [ "$MEM_TOTAL" -lt 1048576 ]; then
        log "⏳ 自动创建 256MB Swap..."
        dd if=/dev/zero of="$USB_MOUNT/swapfile" bs=1M count=256 2>/dev/null
        mkswap "$USB_MOUNT/swapfile" 2>/dev/null
        chmod 600 "$USB_MOUNT/swapfile"
        swapon "$USB_MOUNT/swapfile" 2>/dev/null && log "✅ Swap 已启用"
    fi
elif [ -n "$USB_MOUNT" ] && [ -f "$USB_MOUNT/swapfile" ]; then
    # 原版逻辑：检测到现有 swapfile 就启用
    swapon "$USB_MOUNT/swapfile" 2>/dev/null && log "✅ 已开启现有 Swap"
fi

# ── apk update ───────────────────────────────────────────────────────
apk_update_ok() {
    local out n
    out=$(apk update 2>&1)
    echo "$out" | logger -t "$LOG_TAG"
    n=$(echo "$out" | grep -oE 'OK: [0-9]+ distinct' | grep -oE '[0-9]+')
    [ "${n:-0}" -gt 1000 ] && return 0

    log "⚠️ HTTPS 包数不足（${n:-0}），切换 HTTP..."
    for f in /etc/apk/repositories.d/*.list; do
        [ -f "$f" ] && sed -i \
            's|https://downloads.immortalwrt.org|http://downloads.immortalwrt.org|g' "$f"
    done
    out=$(apk update 2>&1)
    echo "$out" | logger -t "$LOG_TAG"
    n=$(echo "$out" | grep -oE 'OK: [0-9]+ distinct' | grep -oE '[0-9]+')
    [ "${n:-0}" -gt 1000 ] && return 0
    return 1
}

if ! apk_update_ok; then
    log "❌ apk update 失败，等下次重试"
    rm -f "$LOCK"; exit 1
fi
log "✅ apk update 成功"

# ── 安装函数 ─────────────────────────────────────────────────────────
INSTALL_FAIL=0
install_pkg() {
    local out ret
    log "📦 安装: $*"
    out=$(apk add --no-cache "$@" 2>&1)
    ret=$?
    echo "$out" | logger -t "$LOG_TAG"
    if [ "$ret" -eq 0 ]; then
        log "✅ 完成: $*"
    else
        log "⚠️ 失败: $*"
        INSTALL_FAIL=$((INSTALL_FAIL + 1))
    fi
}

# ── 轻量工具 ─────────────────────────────────────────────────────────
install_pkg smartmontools bash htop fdisk lsblk

# ── 中等体积 ─────────────────────────────────────────────────────────
install_pkg samba4-server luci-app-samba4 luci-i18n-samba4-zh-cn

# ── 重型应用 ─────────────────────────────────────────────────────────
if [ -n "$USB_MOUNT" ]; then
    log "⚙️ U 盘已就绪，准备安装重型应用..."

    # Docker
    mkdir -p "$USB_MOUNT/docker-data"
    mkdir -p /etc/docker
    printf '{\n  "data-root": "%s/docker-data",\n  "log-driver": "json-file",\n  "log-opts": {"max-size": "10m", "max-file": "3"}\n}\n' \
        "$USB_MOUNT" > /etc/docker/daemon.json
    log "✅ daemon.json 已写入，data-root → $USB_MOUNT/docker-data"

    install_pkg docker dockerd
    install_pkg luci-app-dockerman luci-i18n-dockerman-zh-cn

    mountpoint -q /sys/fs/cgroup || \
        mount -t cgroup2 none /sys/fs/cgroup 2>/dev/null || true

    /etc/init.d/dockerd enable 2>/dev/null || true
    /etc/init.d/dockerd start  2>/dev/null || true
    log "✅ Docker 已启动，数据目录 → $USB_MOUNT/docker-data"

    # AdGuardHome
    install_pkg adguardhome luci-app-adguardhome

    AGH_USB="$USB_MOUNT/AdGuardHome"
    mkdir -p "$AGH_USB/data"
    mkdir -p "$AGH_USB/config"

    /etc/init.d/adguardhome stop 2>/dev/null || true

    if [ ! -f "$AGH_USB/config/adguardhome.yaml" ] && \
       [ -f "/etc/adguardhome/adguardhome.yaml" ]; then
        cp -a /etc/adguardhome/. "$AGH_USB/config/"
        log "✅ 新环境，初始配置已写入 U 盘 → $AGH_USB/config/"
    else
        log "ℹ️ U 盘已存在旧配置，直接复用，不覆盖"
    fi

    uci set adguardhome.config.config_file="$AGH_USB/config/adguardhome.yaml"
    uci set adguardhome.config.work_dir="$AGH_USB/data"
    uci commit adguardhome
    log "✅ AdGuardHome uci 配置已更新"

    /etc/init.d/adguardhome enable 2>/dev/null || true
    /etc/init.d/adguardhome start  2>/dev/null || true
    log "✅ AdGuardHome 已启动，完全运行于 U 盘"

else
    log "ℹ️ 无 U 盘，跳过 AdGuardHome 和 Docker 安装"
    log "ℹ️ 插入 U 盘后执行: rm /etc/install-extras-done && reboot"
fi

# ── 刷新 LuCI ────────────────────────────────────────────────────────
rm -f /tmp/luci-indexcache* 2>/dev/null
/etc/init.d/rpcd restart 2>/dev/null || true

rm -f /etc/hotplug.d/iface/99-install-extras
rm -f "$LOCK"

# ── 写 DONE_FLAG ─────────────────────────────────────────────────────
if [ "$INSTALL_FAIL" -eq 0 ]; then
    touch "$DONE_FLAG"
    log "🎉 install-extras 全部完成"
else
    log "⚠️ 有 $INSTALL_FAIL 组包失败，下次 WAN UP 重试"
    mkdir -p /etc/hotplug.d/iface
    cat > /etc/hotplug.d/iface/99-install-extras << 'HOTPLUGEOF'
#!/bin/sh
[ "$ACTION" = "ifup" ] || exit 0
[ -f /etc/install-extras-done ] && exit 0
case "$INTERFACE" in
    wan|wan6|pppoe-wan) ;;
    *) exit 0 ;;
esac
/usr/lib/install-extras/run.sh &
HOTPLUGEOF
    chmod +x /etc/hotplug.d/iface/99-install-extras
fi
RUNEOF

chmod +x /usr/lib/install-extras/run.sh

# ── 安装 hotplug 触发器 ──────────────────────────────────────────────
mkdir -p /etc/hotplug.d/iface
cat > /etc/hotplug.d/iface/99-install-extras << 'HOTPLUGEOF'
#!/bin/sh
[ "$ACTION" = "ifup" ] || exit 0
[ -f /etc/install-extras-done ] && exit 0
case "$INTERFACE" in
    wan|wan6|pppoe-wan) ;;
    *) exit 0 ;;
esac
/usr/lib/install-extras/run.sh &
HOTPLUGEOF

chmod +x /etc/hotplug.d/iface/99-install-extras

logger -t "96-install-extras" "✅ hotplug 触发器已安装（优化版）"
exit 0
