# files/etc/uci-defaults/96-install-extras
#!/bin/sh
# 96-install-extras [v2 - Fusion Edition]
# 融合：v5.1 架构（离线安装/Flash 友好）+ v1 优化（简洁/严谨/灵活）
# 修复：/tmp 空间不足 / Swap 持久化 / ash 兼容 / 竞态防护
# 适用：360 V6 / 小米 CR660x / IPQ6000 (128M Flash / ImmortalWrt)

DONE_FLAG="/etc/install-extras-done"
USB_DONE="/etc/install-usb-done"
STATE_DIR="/etc/install-state"

mkdir -p "$STATE_DIR" /usr/lib/install-extras /usr/bin \
         /etc/hotplug.d/iface /etc/hotplug.d/block \
         /etc/dbus-1 /etc/profile.d

# ====================================================================
# 兼容函数（ash 无 mountpoint 命令）
# ====================================================================
_is_mounted() {
    grep -q "^[^ ]* $1 " /proc/mounts 2>/dev/null
}

# ====================================================================
# 模块 0：360v6-status 状态查询工具
# ====================================================================
cat > /usr/bin/360v6-status << 'STATUSEOF'
#!/bin/sh
STATE_DIR="/etc/install-state"
_is_mounted() { grep -q "^[^ ]* $1 " /proc/mounts 2>/dev/null; }

echo "=== 360 V6 部署状态 @ $(date '+%Y-%m-%d %H:%M:%S') ==="
echo ""

[ -f /etc/install-extras-done ] && echo "[OK] 基础环境：已完成" || echo "[..] 基础环境：待触发"
p1=$(cat "$STATE_DIR/phase1-failures" 2>/dev/null || echo "0")
[ "$p1" -gt 0 ] && echo "     重试: $p1/3"

[ -f /etc/install-usb-done ] && echo "[OK] 重型应用：已完成" || echo "[..] 重型应用：待 U 盘"
p2=$(cat "$STATE_DIR/phase2-failures" 2>/dev/null || echo "0")
[ "$p2" -gt 0 ] && echo "     重试: $p2/3"

echo ""
df -h /overlay 2>/dev/null | awk 'NR==2 {print "Flash : "$3"/"$2" ("$5")"}'
free -h 2>/dev/null | awk '/Mem:/ {print "Mem   : "$3"/"$2}'
echo ""

for mp in /mnt/sda1 /mnt/sda /mnt/sdb1 /mnt/sdb /mnt/usb /mnt/mmcblk0p1; do
    _is_mounted "$mp" && { echo "USB   : $mp"; break; }
done
grep -q swapfile /proc/swaps 2>/dev/null && echo "Swap  : ✅ 已激活"

docker info >/dev/null 2>&1 && echo "[OK] Docker" || echo "[--] Docker"
pgrep AdGuardHome >/dev/null 2>&1 && echo "[OK] AGH" || echo "[--] AGH"
echo ""
echo "重跑: rm -f /etc/install-*-done $STATE_DIR/*-failures && /usr/lib/install-extras/run.sh"
STATUSEOF
chmod +x /usr/bin/360v6-status

# ====================================================================
# 模块 A：install-usb.sh（Phase 2 重型应用）✅ v2 融合版
# ====================================================================
cat > /usr/lib/install-extras/install-usb.sh << 'USBEOF'
#!/bin/sh
# ✅ v2 融合：v5.1 架构 + v1 优化（简洁/严谨/灵活）
LOG_TAG="usb-extras"
USB_DONE="/etc/install-usb-done"
USB_LOCK="/tmp/install-usb.lock"
STATE_DIR="/etc/install-state"

log() { echo "[$(date '+%H:%M:%S')] $1" | logger -t "$LOG_TAG"; }

mkdir -p "$STATE_DIR"
[ -f "$USB_DONE" ] && { log "ℹ️ 已完成，跳过"; exit 0; }
[ -f "$USB_LOCK" ] && { log "⚠️ 运行中，跳过"; exit 0; }

__fail=$(cat "$STATE_DIR/phase2-failures" 2>/dev/null || echo "0")
[ "${__fail:-0}" -ge 3 ] && { log "❌ 失败过多"; exit 1; }
touch "$USB_LOCK"

# ── DNS 恢复 ─────────────────────────────────────────────────────
__dns_restore=0
__dns_old=""
__restore_dns() {
    [ "$__dns_restore" = "1" ] || return
    printf '%s' "${__dns_old:-nameserver 127.0.0.1}" > /etc/resolv.conf
}
trap 'rm -f "$USB_LOCK"; __restore_dns' EXIT INT TERM

# ── 确定挂载点（ash 兼容）────────────────────────────────────────
USB_MOUNT="${1:-}"
if [ -z "$USB_MOUNT" ]; then
    for __mp in /mnt/sda1 /mnt/sda /mnt/sdb1 /mnt/sdb /mnt/usb /mnt/mmcblk0p1; do
        _is_mounted "$__mp" && { USB_MOUNT="$__mp"; break; }
    done
fi
[ -z "$USB_MOUNT" ] && USB_MOUNT=$(grep -E '^/dev/sd[a-z][0-9]* ' /proc/mounts 2>/dev/null | head -1 | awk '{print $2}')
[ -z "$USB_MOUNT" ] && { log "❌ 未找到 U 盘"; exit 1; }
log "▶ 重型应用部署 | $USB_MOUNT"

# ── 记录挂载点（供 init 脚本读取）────────────────────────────────
echo "$USB_MOUNT" > /etc/install-usb-mount

# ── ✅ 关键修复：重定向 apk 缓存 + TMPDIR 到 U 盘（解决 /tmp 爆满）──
APK_CACHE_DIR="$USB_MOUNT/.apk-cache"
TMPDIR="$USB_MOUNT/.tmp"
mkdir -p "$APK_CACHE_DIR" "$TMPDIR"
export TMPDIR
log "📁 缓存: $APK_CACHE_DIR"

# ── U 盘空间预检 ─────────────────────────────────────────────────
[ "$(df -k "$USB_MOUNT" 2>/dev/null | awk 'NR==2{print $4}')" -lt 2097152 ] && {
    log "⚠️ U 盘空间不足"; touch "$USB_DONE"; exit 0; }

# ── DNS 保护 ─────────────────────────────────────────────────────
if ! nslookup baidu.com >/dev/null 2>&1; then
    __dns_old=$(cat /etc/resolv.conf 2>/dev/null)
    printf 'nameserver 223.5.5.5\nnameserver 8.8.8.8\n' > /etc/resolv.conf
    __dns_restore=1
fi

# ── I/O 调度 + 内存优化 ──────────────────────────────────────────
USB_DEV=$(grep " $USB_MOUNT " /proc/mounts 2>/dev/null | head -1 | awk '{print $1}' | sed 's|/dev/||;s|[0-9]*$||')
[ -f "/sys/block/$USB_DEV/queue/scheduler" ] && {
    grep -q "mq-deadline" "/sys/block/$USB_DEV/queue/scheduler" 2>/dev/null && echo "mq-deadline" > "/sys/block/$USB_DEV/queue/scheduler" 2>/dev/null || echo "deadline" > "/sys/block/$USB_DEV/queue/scheduler" 2>/dev/null
}
echo "10" > /proc/sys/vm/swappiness 2>/dev/null
echo "16384" > /proc/sys/vm/min_free_kbytes 2>/dev/null

# ── ✅ v1 优化：Swap 激活成功后才写 fstab（严谨持久化）────────────
SWAPFILE="$USB_MOUNT/swapfile"
if ! grep -q swapfile /proc/swaps 2>/dev/null; then
    [ ! -f "$SWAPFILE" ] && {
        log "⏳ 创建 Swap 256MB..."
        dd if=/dev/zero of="${SWAPFILE}.tmp" bs=1M count=256 2>/dev/null && {
            mv "${SWAPFILE}.tmp" "$SWAPFILE"; sync; chmod 600 "$SWAPFILE"; mkswap "$SWAPFILE" >/dev/null 2>&1
            log "✅ Swap 文件已创建"
        } || rm -f "${SWAPFILE}.tmp"
    }
    if [ -f "$SWAPFILE" ] && swapon "$SWAPFILE" 2>/dev/null; then
        log "✅ Swap 已激活"
        grep -qF "$SWAPFILE" /etc/fstab 2>/dev/null || echo "$SWAPFILE none swap sw 0 0" >> /etc/fstab
        log "✅ Swap 已写入 fstab"
    fi
fi

# ── 等待 apk 就绪 ────────────────────────────────────────────────
for __i in $(seq 1 30); do apk info >/dev/null 2>&1 && break; sleep 2; done
apk info >/dev/null 2>&1 || { log "❌ apk 未就绪"; exit 1; }

# ── 安装函数（简洁变量命名）──────────────────────────────────────
__fail_cnt=0
__install() {
    __tmp=$(mktemp "$TMPDIR/apk.XXXXXX")
    log "📦 $*"
    apk add --no-cache "$@" > "$__tmp" 2>&1
    __ret=$?
    cat "$__tmp" | logger -t "$LOG_TAG"
    rm -f "$__tmp"
    [ "$__ret" -eq 0 ] && log "✅ $*" || { log "❌ $*"; __fail_cnt=$((__fail_cnt + 1)); }
    return $__ret
}

# ====================================================================
# Docker 安装到 U 盘（v5.1 架构 + v1 简洁）
# ====================================================================
log "📦 Docker 部署"
DOCKER_SYS="$USB_MOUNT/docker-system"
DOCKER_DATA="$USB_MOUNT/docker-data"
mkdir -p "$DOCKER_SYS" "$DOCKER_DATA" /etc/docker

# 初始化 U 盘 apk 环境
mkdir -p "$DOCKER_SYS/etc/apk" "$DOCKER_SYS/var/cache/apk"
cp -f /etc/apk/repositories "$DOCKER_SYS/etc/apk/" 2>/dev/null
[ ! -f "$DOCKER_SYS/etc/apk/world" ] && apk add --root "$DOCKER_SYS" --initdb --no-cache 2>/dev/null

# 安装 Docker 到 U 盘（✅ 使用重定向缓存）
__install --root "$DOCKER_SYS" --cache-dir "$APK_CACHE_DIR" docker dockerd 2>/dev/null

__docker_ok=0
if [ -f "$DOCKER_SYS/usr/bin/dockerd" ]; then
    # ✅ 符号链接写到 overlay 可写层（重启不丢失）
    mkdir -p /overlay/upper/usr/bin
    for __bin in docker dockerd containerd runc; do
        ln -sf "$DOCKER_SYS/usr/bin/$__bin" "/overlay/upper/usr/bin/$__bin" 2>/dev/null
    done
    # LD_LIBRARY_PATH 持久化
    echo "export LD_LIBRARY_PATH=\"$DOCKER_SYS/usr/lib:\$LD_LIBRARY_PATH\"" > /etc/profile.d/docker-usb.sh

    # daemon.json（路径写死）
    cat > /etc/docker/daemon.json << EOF
{"data-root":"$DOCKER_DATA","storage-driver":"overlay2","log-driver":"json-file","log-opts":{"max-size":"5m","max-file":"2"}}
EOF

    # ✅ init 脚本读取固定配置（不运行时扫描）
    cat > /etc/init.d/dockerd << 'INITEOF'
#!/bin/sh /etc/rc.common
START=99; STOP=10
start() {
    [ -f /etc/install-usb-done ] || return 0
    USB_MOUNT=$(cat /etc/install-usb-mount 2>/dev/null)
    [ -f "$USB_MOUNT/docker-system/usr/bin/dockerd" ] || return 1
    export LD_LIBRARY_PATH="$USB_MOUNT/docker-system/usr/lib:$LD_LIBRARY_PATH"
    $USB_MOUNT/docker-system/usr/bin/dockerd --config-file /etc/docker/daemon.json &
}
stop() { pkill dockerd 2>/dev/null; }
INITEOF
    chmod +x /etc/init.d/dockerd
    /etc/init.d/dockerd enable >/dev/null 2>&1
    /etc/init.d/dockerd start >/dev/null 2>&1
    for __i in $(seq 1 15); do docker info >/dev/null 2>&1 && { __docker_ok=1; break; }; sleep 2; done
    [ "$__docker_ok" -eq 1 ] && log "✅ Docker 已启动" || log "⚠️ Docker 启动超时"
else
    log "❌ Docker 安装失败"
fi

# ====================================================================
# AdGuardHome 安装到 U 盘
# ====================================================================
log "📦 AdGuardHome 部署"
AGH_SYS="$USB_MOUNT/adguardhome-system"
AGH_DATA="$USB_MOUNT/adguardhome-data"
AGH_CONF="$USB_MOUNT/adguardhome-config"
mkdir -p "$AGH_SYS" "$AGH_DATA" "$AGH_CONF"

[ ! -f "$AGH_SYS/etc/apk/world" ] && apk add --root "$AGH_SYS" --initdb --no-cache 2>/dev/null
__install --root "$AGH_SYS" --cache-dir "$APK_CACHE_DIR" adguardhome 2>/dev/null

__agh_ok=0
if [ -f "$AGH_SYS/usr/bin/AdGuardHome" ]; then
    ln -sf "$AGH_SYS/usr/bin/AdGuardHome" "/overlay/upper/usr/bin/AdGuardHome" 2>/dev/null

    # 配置迁移
    [ ! -f "$AGH_CONF/AdGuardHome.yaml" ] && [ -f "/etc/adguardhome/AdGuardHome.yaml" ] && {
        cp -a /etc/adguardhome/. "$AGH_CONF/" 2>/dev/null; log "✅ AGH 配置已迁移"; }
    [ -z "$(ls -A "$AGH_DATA" 2>/dev/null)" ] && [ -d "/var/lib/adguardhome" ] && {
        cp -a /var/lib/adguardhome/. "$AGH_DATA/" 2>/dev/null; log "✅ AGH 数据已迁移"; }

    uci set adguardhome.config.config_file="$AGH_CONF/AdGuardHome.yaml" 2>/dev/null
    uci set adguardhome.config.work_dir="$AGH_DATA" 2>/dev/null
    uci commit adguardhome 2>/dev/null

    # ✅ init 脚本读取固定路径
    cat > /etc/init.d/adguardhome << 'AGHINIT'
#!/bin/sh /etc/rc.common
START=98; STOP=10
start() {
    [ -f /etc/install-usb-done ] || return 0
    USB_MOUNT=$(cat /etc/install-usb-mount 2>/dev/null)
    [ -f "$USB_MOUNT/adguardhome-system/usr/bin/AdGuardHome" ] || return 1
    $USB_MOUNT/adguardhome-system/usr/bin/AdGuardHome --config "$USB_MOUNT/adguardhome-config/AdGuardHome.yaml" --work-dir "$USB_MOUNT/adguardhome-data" &
}
stop() { pkill AdGuardHome 2>/dev/null; }
AGHINIT
    chmod +x /etc/init.d/adguardhome
    /etc/init.d/adguardhome enable >/dev/null 2>&1
    /etc/init.d/adguardhome restart >/dev/null 2>&1
    sleep 3
    pgrep -x AdGuardHome >/dev/null 2>&1 && __agh_ok=1
    [ "$__agh_ok" -eq 1 ] && log "✅ AGH 已启动" || log "⚠️ AGH 进程未检测到"
else
    log "❌ AGH 安装失败"
fi

# ── 刷新 LuCI ────────────────────────────────────────────────────
rm -f /tmp/luci-indexcache* 2>/dev/null
/etc/init.d/rpcd restart >/dev/null 2>&1

# ── ✅ v1 优化：部分成功即标记完成（灵活决策）────────────────────
[ "$__docker_ok" -eq 1 ] || [ "$__agh_ok" -eq 1 ] && {
    touch "$USB_DONE"
    rm -f "$STATE_DIR/phase2-failures"
    log "🎉 完成 [Docker:$([ "$__docker_ok" -eq 1 ]&&echo OK||echo FAIL) AGH:$([ "$__agh_ok" -eq 1 ]&&echo OK||echo FAIL)]"
    exit 0
}
echo $((__fail + 1)) > "$STATE_DIR/phase2-failures"
log "❌ 均失败 (重试 $((__fail + 1))/3)"
exit 1
USBEOF
chmod +x /usr/lib/install-extras/install-usb.sh

# ====================================================================
# 模块 B：run.sh（Phase 1 基础环境）
# ====================================================================
cat > /usr/lib/install-extras/run.sh << 'RUNEOF'
#!/bin/sh
LOG_TAG="install-main"
DONE_FLAG="/etc/install-extras-done"
USB_DONE="/etc/install-usb-done"
RUN_LOCK="/tmp/install-main.lock"
STATE_DIR="/etc/install-state"

log() { echo "[$(date '+%H:%M:%S')] $1" | logger -t "$LOG_TAG"; }

mkdir -p "$STATE_DIR"
[ -f "$DONE_FLAG" ] && { log "ℹ️ 已完成"; exit 0; }
[ -f "$RUN_LOCK" ] && { log "⚠️ 运行中"; exit 0; }

__fail=$(cat "$STATE_DIR/phase1-failures" 2>/dev/null || echo "0")
[ "${__fail:-0}" -ge 3 ] && { log "❌ 失败过多"; exit 1; }
touch "$RUN_LOCK"

__dns_restore=0
__dns_old=""
__restore_dns() {
    rm -f "$RUN_LOCK"
    [ "$__dns_restore" = "1" ] || return
    printf '%s' "${__dns_old:-nameserver 127.0.0.1}" > /etc/resolv.conf
    /etc/init.d/dnsmasq restart >/dev/null 2>&1
}
trap '__restore_dns' EXIT INT TERM

log "▶ v2 基础环境启动"

# 清理无效源
for __f in /etc/apk/repositories.d/*.list; do [ -f "$__f" ] && sed '/\/video\//d' "$__f" > "${__f}.tmp" && mv "${__f}.tmp" "$__f"; done 2>/dev/null

# 等待网络
sleep 10
for __i in $(seq 1 36); do ping -c1 -W2 223.5.5.5 >/dev/null 2>&1 && break; ping -c1 -W2 8.8.8.8 >/dev/null 2>&1 && break; sleep 5; done
ping -c1 -W2 223.5.5.5 >/dev/null 2>&1 || { log "❌ 网络不可达"; echo $((__fail + 1)) > "$STATE_DIR/phase1-failures"; exit 1; }

# DNS 注入
__dns_old=$(cat /etc/resolv.conf 2>/dev/null)
printf 'nameserver 223.5.5.5\nnameserver 8.8.8.8\n' > /etc/resolv.conf
__dns_restore=1

# Flash 预检
FLASH_FREE=$(df /overlay 2>/dev/null | awk 'NR==2 {print $4}')
[ "${FLASH_FREE:-0}" -lt 20480 ] && { log "❌ Flash 不足 20MB"; exit 1; }

# apk update
__apk_ok() {
    __tmp=$(mktemp /tmp/apk.XXXXXX)
    apk update > "$__tmp" 2>&1
    __n=$(grep -oE 'OK: [0-9]+ distinct' "$__tmp" 2>/dev/null | grep -oE '[0-9]+' | head -1)
    rm -f "$__tmp"
    [ "${__n:-0}" -gt 800 ] && return 0
    for __f in /etc/apk/repositories.d/*.list; do [ -f "$__f" ] && sed 's|https://|http://|g' "$__f" > "${__f}.tmp" && mv "${__f}.tmp" "$__f"; done 2>/dev/null
    __tmp=$(mktemp /tmp/apk.XXXXXX)
    apk update > "$__tmp" 2>&1
    __n=$(grep -oE 'OK: [0-9]+ distinct' "$__tmp" 2>/dev/null | grep -oE '[0-9]+' | head -1)
    rm -f "$__tmp"
    [ "${__n:-0}" -gt 800 ]
}
__apk_ok || { log "❌ apk update 失败"; exit 1; }

# 安装函数
__install_fail=0
__install() {
    __tmp=$(mktemp /tmp/apk.XXXXXX)
    log "📦 $*"
    apk add --no-cache "$@" > "$__tmp" 2>&1
    __ret=$?
    cat "$__tmp" | logger -t "$LOG_TAG"
    rm -f "$__tmp"
    [ "$__ret" -eq 0 ] && log "✅ $*" || { log "❌ $*"; __install_fail=$((__install_fail + 1)); }
}

# 基础工具 + Samba（25MB 阈值）
for __pkg in bash htop lsblk smartmontools fdisk; do __install "$__pkg"; done
[ "$FLASH_FREE" -ge 25000 ] && __install samba4-server luci-app-samba4 luci-i18n-samba4-zh-cn || log "ℹ️ 跳过 Samba"

# 修复 avahi
/etc/init.d/avahi-daemon stop >/dev/null 2>&1; sleep 1
/etc/init.d/dbus stop >/dev/null 2>&1; sleep 1
/etc/init.d/dbus start >/dev/null 2>&1; sleep 2
/etc/init.d/avahi-daemon start >/dev/null 2>&1

# 收尾
[ "$__install_fail" -eq 0 ] && { touch "$DONE_FLAG"; rm -f "$STATE_DIR/phase1-failures"; log "🎉 基础环境完成"; } || { echo $((__fail + 1)) > "$STATE_DIR/phase1-failures"; exit 1; }

# U 盘检测（单次检查 + 后台触发）
[ ! -f "$USB_DONE" ] && {
    USB_MOUNT=""
    for __mp in /mnt/sda1 /mnt/sda /mnt/sdb1 /mnt/sdb /mnt/usb /mnt/mmcblk0p1; do _is_mounted "$__mp" && { USB_MOUNT="$__mp"; break; }; done
    [ -n "$USB_MOUNT" ] && { log "✅ U 盘: $USB_MOUNT"; /usr/lib/install-extras/install-usb.sh "$USB_MOUNT" &; } || log "ℹ️ 等待 U 盘热插拔"
}
log "💡 执行 '360v6-status' 查看状态"
RUNEOF
chmod +x /usr/lib/install-extras/run.sh

# ====================================================================
# 模块 C/D：hotplug 触发器
# ====================================================================
cat > /etc/hotplug.d/iface/96-install-extras << 'IFACEEOF'
#!/bin/sh
[ "${ACTION:-}" = "ifup" ] || exit 0
case "${INTERFACE:-}" in wan|wan6|pppoe-wan) ;; *) exit 0 ;; esac
[ -f /etc/install-extras-done ] && exit 0
logger -t "hotplug" "🌐 WAN 上线 → 基础环境"
/usr/lib/install-extras/run.sh &
IFACEEOF
chmod +x /etc/hotplug.d/iface/96-install-extras

cat > /etc/hotplug.d/block/96-install-usb << 'BLOCKEOF'
#!/bin/sh
[ "${ACTION:-}" = "add" ] || exit 0
[ -f /etc/install-usb-done ] && exit 0
[ -f /etc/install-extras-done ] || exit 0
logger -t "hotplug" "💾 U 盘插入: ${DEVNAME:-?}"
for __i in $(seq 1 15); do
    for __mp in /mnt/sda1 /mnt/sda /mnt/sdb1 /mnt/sdb /mnt/usb /mnt/mmcblk0p1; do
        _is_mounted "$__mp" && { logger -t "hotplug" "✅ 挂载: $__mp"; /usr/lib/install-extras/install-usb.sh "$__mp" &; exit 0; }
    done
    sleep 2
done
logger -t "hotplug" "⚠️ 挂载点未就绪"
BLOCKEOF
chmod +x /etc/hotplug.d/block/96-install-usb

# ====================================================================
# 完成
# ====================================================================
logger -t "install-extras" "✅ v2 Fusion 就绪 | Flash 友好 + /tmp 修复"
echo "━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━"
echo "🚀 v2 Fusion Edition 部署完毕"
echo "💡 监控: logread -f | grep -E 'usb-extras|✅|❌'"
echo "💡 状态: 360v6-status"
echo "━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━"
exit 0
