# files/etc/uci-defaults/97-fix-apk-feeds
#!/bin/sh
# 97-fix-apk-feeds
# 修复 video feed，使用安全的 sed 语法

LOG_TAG="fix-apk-feeds"

log() {
    echo "[$(date '+%H:%M:%S')] $1" | logger -t "$LOG_TAG"
}

REPO_LIST_DIR="/etc/apk/repositories.d"
FIXED_COUNT=0

if [ ! -d "$REPO_LIST_DIR" ]; then
    log "⚠️  $REPO_LIST_DIR 不存在，跳过修复"
    exit 0
fi

for f in "$REPO_LIST_DIR"/*.list; do
    [ -f "$f" ] || continue
    
    # 使用 | 作为分隔符，避免路径中的 / 冲突
    if sed -i '\|/video/|d' "$f" 2>/dev/null; then
        FIXED_COUNT=$((FIXED_COUNT + 1))
        log "修复: $(basename "$f")"
    else
        log "⚠️  修复失败: $(basename "$f")"
    fi
done

if [ $FIXED_COUNT -gt 0 ]; then
    log "✅ 修复了 $FIXED_COUNT 个 feeds 文件"
else
    log "ℹ️  没有需要修复的 feeds"
fi

exit 0
