# files/etc/init.d/bbr-load
#!/bin/sh /etc/rc.common
# 加载 BBR 拥塞控制并应用 sysctl 参数
START=15
STOP=89

start() {
    modprobe tcp_bbr 2>/dev/null || true
    echo 1 > /proc/sys/net/ipv4/ip_forward 2>/dev/null || true
    sysctl -p /etc/sysctl.d/99-bbr.conf 2>/dev/null || true
    logger -t bbr-load "BBR loaded, sysctl applied"
}

stop() {
    logger -t bbr-load "bbr-load stop (no-op)"
}
